# etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
/opt/rh/rh-nginx120/root/usr/sbin/nginx -g "daemon off;"
